# ubuntu/ghost.upstart
description "dojo"
# ghost - nodejs service wrapper
#
# Just a blogging platform.

start on (filesystem and net-device-up)
stop on runlevel [!2345]

env RUN_AS="ghost"
env APP="/opt/ghost/index.js"
env LOG="/var/log/ghost/ghost.log"

# app settings
env NODE_ENV="production"

# If the process quits unexpectadly trigger a respawn
respawn

# Start the process
exec start-stop-daemon --start --chuid $RUN_AS --make-pidfile --pidfile /var/run/ghost.pid --exec /usr/bin/node -- $APP >> $LOG 2>&1
